# conf.d/supervisor/supervisor.conf
[supervisord]
nodaemon = true

[program:gunicorn-app]
command = /usr/local/bin/gunicorn --chdir /git-repos/myweb/src conf.wsgi:application --bind=127.0.0.1:8020 --workers 1 --log-level debug
redirect_stderr=true

[program:nginx-app]
command = /usr/sbin/nginx